# conf/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1920x1080x24 +extension RANDR
autorestart=true
priority=100

[program:chrome]
environment=HOME="/home/chrome",DISPLAY=":1",USER="chrome"
command=/opt/google/chrome/chrome --user-data-dir  --no-sandbox --window-position=0,0 --window-size=1920,1080 --force-device-scale-factor=1 https://www.google.com
user=chrome
autorestart=true
priority=200

[program:x11vnc]
command=/usr/bin/x11vnc -display :1 %(ENV_X11VNC_AUTH)s -wait 5 -forever -xrandr
user=chrome
autorestart=true
priority=300

[program:xrdp-sesman]
command=/usr/sbin/xrdp-sesman --nodaemon
user=root
process_name = xrdp-sesman
autorestart=true
priority=400

[program:xrdp]
##command=/usr/sbin/xrdp -nodaemon
command=/etc/init.d/xrdp restart
user=root
process_name = xrdp
autorestart=true
stderr_logfile=/var/log/xrdp-error.log
stdout_logfile=/var/log/xrdp.log
priority=500
